# systemd/sagequeue-container.service
[Unit]
Description=SageQueue: ensure podman-compose container is up

[Service]
Type=oneshot
EnvironmentFile=%h/.config/sagequeue/sagequeue.env
Environment=PATH=/usr/bin:/bin
ExecStart=/bin/bash -c 'set -euo pipefail; exec "$PROJECT_ROOT/bin/sagequeue-ensure-container.sh"'
RemainAfterExit=yes

[Install]
WantedBy=default.target

