# systemd/sagequeue-container.service
[Unit]
Description=SageQueue: ensure podman-compose container is up

[Service]
Type=oneshot
EnvironmentFile=%h/.config/sagequeue/sagequeue.env
ExecStart=/usr/bin/env bash -c '. %h/.config/sagequeue/sagequeue.env; exec "$PROJECT_ROOT/bin/sagequeue-ensure-container.sh"'
RemainAfterExit=yes

[Install]
WantedBy=default.target

