# systemd/sagequeue-recover.service
[Unit]
Description=SageQueue: requeue orphaned running jobs

[Service]
Type=oneshot
EnvironmentFile=%h/.config/sagequeue/sagequeue.env
ExecStart=/usr/bin/env bash -lc '. %h/.config/sagequeue/sagequeue.env; exec "$PROJECT_ROOT/bin/sagequeue-recover.sh"'

