# systemd/sagequeue@.service
[Unit]
Description=SageQueue worker %i
After=sagequeue-container.service
Requires=sagequeue-container.service

[Service]
Type=simple
EnvironmentFile=%h/.config/sagequeue/sagequeue.env
Environment=PATH=/usr/bin:/bin

# Requeue orphans before starting (safe: recover uses flock)
ExecStartPre=/bin/bash -c 'set -euo pipefail; exec "$PROJECT_ROOT/bin/sagequeue-recover.sh"'
ExecStart=/bin/bash -c 'set -euo pipefail; exec "$PROJECT_ROOT/bin/sagequeue-worker.sh" %i'
Restart=always
RestartSec=2
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target

